# r/supervisord.conf
[program:r-studio-server]
process_name=rstudio
directory=/data
command=/bin/bash /root/start_rstudio_server.sh
stdout_logfile=/var/log/r_studio_server.log
redirect_stderr=true
priority=900
